# files/apache.cron
#!/bin/bash
/usr/bin/yum -y update httpd
