# files/apache.cron
#!/bin/bash
/usr/bin/yum -y update httpd
/etc/init.d/httpd reload
